# FILES/queue.conf
[program:queue]
command=php artisan queue:listen --tries=2 --timeout=600
directory=/vagrant/cnp
stdout_logfile=/vagrant/cnp/app/storage/logs/supervisor.log
redirect_stderr=true

